# root/etc/services.d/geph4-exit/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid abc /usr/bin/geph4-exit --config /geph4-exit.toml
